# bilibili-api.service
[Unit]
Description=Bilibili Toolkit API Server
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/bilibili-toolkit
Environment=BILI_TOOLKIT_TOKEN=your_token_here
Environment=BILI_TOOLKIT_PORT=18020
Environment=BILI_TOOLKIT_CONFIG=/opt/bilibili-toolkit/config.json
Environment=BILI_TOOLKIT_SESSDATA=your_sessdata_here
Environment=BILI_TOOLKIT_BILI_JCT=your_bili_jct_here
Environment=BILI_TOOLKIT_BUVID3=your_buvid3_here
ExecStart=/opt/bilibili-toolkit/venv/bin/python3 /opt/bilibili-toolkit/api_server.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
